# packaging/deb/systemd/st2timersengine.service
[Unit]
Description=StackStorm service st2timersengine
After=network.target

[Service]
Type=simple
User=st2
Group=st2
Environment="DAEMON_ARGS=--config-file /etc/st2/st2.conf"
EnvironmentFile=-/etc/default/st2timersengine
ExecStart=/opt/stackstorm/st2/bin/st2timersengine $DAEMON_ARGS
TimeoutSec=60
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
