# packages/st2/debian/st2api.socket
[Unit]
Description=StackStorm st2api Socket
PartOf=st2api.service

[Socket]
ListenStream=127.0.0.1:9101

[Install]
WantedBy=sockets.target
