# packages/st2api/rpm/st2api.service
[Unit]
Description=StackStorm service st2api
After=network.target

[Service]
Type=simple
User=st2
Group=st2
EnvironmentFile=-/etc/sysconfig/st2api
ExecStart=/usr/share/python/st2api/bin/gunicorn_pecan /usr/share/python/st2api/lib/python2.7/site-packages/st2api/gunicorn_config.py -k eventlet -b 127.0.0.1:9101 --threads 10 --workers 1
TimeoutSec=60
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
