# packages/st2auth/rpm/st2auth.service
[Unit]
Description=StackStorm service st2auth
After=network.target

[Service]
Type=simple
User=st2
Group=st2
EnvironmentFile=-/etc/sysconfig/st2auth
ExecStart=/usr/share/python/st2auth/bin/uwsgi --ini /etc/uwsgi.d/st2auth.ini
TimeoutSec=60
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
