# packages/st2auth/rpm/st2auth.service
[Unit]
Description=StackStorm service st2auth
After=network.target

[Service]
Type=simple
User=st2
Group=st2
EnvironmentFile=-/etc/sysconfig/st2auth
ExecStart=/usr/share/python/st2auth/bin/gunicorn_pecan /usr/share/python/st2auth/lib/python2.7/site-packages/st2auth/gunicorn_config.py -k eventlet -b 127.0.0.1:9100 --threads 10 --workers 1
TimeoutSec=60
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
